# archived/capistrano/bungeni_setup/branches/bungeni-10.04/config/templates/supervisord.conf.erb
[inet_http_server]
port = 8888
username = admin
password = admin

[supervisord]
logfile = <%=buildout_dir%>/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = <%=buildout_dir%>/supervisord.pid
nodaemon = false

[supervisorctl]
serverurl = http://<%=app_host%>:8888
username = admin
password = admin

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface


[program:bungeni]
command = <%=buildout_dir%>/bin/paster serve <%=buildout_dir%>/deploy.ini --log-file=<%=buildout_dir%>/bungeni.log --pid-file=<%=buildout_dir%>/bungeni.pid
process_name = bungeni
directory = <%=buildout_dir%>/bin
priority = 10
redirect_stderr = true

[program:plone]
command = <%=plone_buildout_dir%>/bin/paster serve <%=plone_buildout_dir%>/etc/deploy.ini --log-file=<%=plone_buildout_dir%>/plone.log --pid-file=<%=plone_buildout_dir%>/plone.pid
process_name = plone
directory = <%=plone_buildout_dir%>/bin
priority = 20
redirect_stderr = true

[program:portal]
command = <%=portal_buildout_dir%>/bin/paster serve <%=portal_buildout_dir%>/deploy.ini --log-file=<%=portal_buildout_dir%>/portal.log --pid-file=<%=portal_buildout_dir%>/portal.pid
process_name = portal
directory = <%=portal_buildout_dir%>/bin
priority = 30
redirect_stderr = true

[program:postgres]
command = <%=buildout_dir%>/parts/postgresql/bin/postmaster -D <%=buildout_dir%>/parts/db
stdout_logfile = <%=buildout_dir%>/log/postgres.log
priority = 40
redirect_stderr = true
stopsignal = INT


[program:openoffice]
command = soffice -accept='socket,port=2002;urp;' -nologo -nofirststartwizard -headless
stdout_logfile= <%=buildout_dir%>/openoffice.log
priority = 70
redirect_stderr = true
;stopsignal = INT


;
;
;command = <%=buildout_dir%>/bin/pg_ctl start
;process_name = portal
;directory = <%=portal_buildout_dir%>/bin
;priority = 10
;redirect_stderr = true
;


